# roles/nixconreg/files/crontab
SHELL=/bin/bash
5 * * * * cd /var/www/nixconreg && /var/www/nixconreg/venv/bin/python3 manage.py generate_participants_json /var/www/nixconreg/registration/static/participants.json

